# itamae/cookbooks/base/files/etc/profile.d/login_info.sh
echo "---
AWS Account ID:    $AWS_ACCOUNT_ID
Instance ID:       $INSTANCE_ID
Instance Type:     $INSTANCE_TYPE
Availability Zone: $AVAILABILITY_ZONE
Public DNS Name:   $PUBLIC_HOSTNAME
"
